# templates/x708-bat.service.j2
[Unit]
Description=Low battery monitor

[Service]
Type=simple
ExecStart={{ x708_install_path }}/bat.py
Restart=always

User=root
Group=root

[Install]
WantedBy=multi-user.target
